# packages/voice-mode/voice_mode-7.4.0-py3-none-any.whl/voice_mode/templates/launchd/com.voicemode.livekit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- com.voicemode.livekit.plist v1.1.0 -->
<!-- Last updated: 2025-12-01 -->
<!-- Compatible with: LiveKit v1.7.0+ -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.voicemode.livekit</string>
    <key>ProgramArguments</key>
    <array>
        <string>{LIVEKIT_BIN}</string>
        <string>--config</string>
        <string>{HOME}/.voicemode/config/livekit.yaml</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{HOME}/.voicemode/logs/livekit/livekit.out.log</string>
    <key>StandardErrorPath</key>
    <string>{HOME}/.voicemode/logs/livekit/livekit.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{HOME}/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>