# packages/freva-deployment/freva_deployment-2601.0.0.tar.gz/freva_deployment-2601.0.0/assets/share/freva/deployment/playbooks/templates/reverse-proxy-nginx.service.j2
[Unit]
Description=NGINX reverse proxy for {{ service }}
After=network.target

[Service]
ExecStart=/usr/sbin/nginx -g 'daemon off;' -c  {{ nginx_config_path }}
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
Restart=always
RestartSec=5
TimeoutStopSec=10
KillMode=process

[Install]
WantedBy=multi-user.target
