# packages/voice-mode/voice_mode-7.3.0-py3-none-any.whl/voice_mode/templates/launchd/com.voicemode.whisper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- com.voicemode.whisper.plist v1.2.0 -->
<!-- Last updated: 2025-12-01 -->
<!-- Simplified: start script handles config via voicemode.env -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.voicemode.whisper</string>
    <key>ProgramArguments</key>
    <array>
        <string>{START_SCRIPT}</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{HOME}/.voicemode/logs/whisper/whisper.out.log</string>
    <key>StandardErrorPath</key>
    <string>{HOME}/.voicemode/logs/whisper/whisper.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{HOME}/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>