# packages/voice-mode/voice_mode-7.2.0.tar.gz/voice_mode-7.2.0/voice_mode/templates/systemd/voicemode-whisper.service
# voicemode-whisper.service v1.2.0
# Last updated: 2025-12-01
# Simplified: start script handles config via voicemode.env

[Unit]
Description=Whisper.cpp Speech Recognition Server
After=network.target

[Service]
Type=simple
ExecStart={START_SCRIPT}
Restart=on-failure
RestartSec=10

# Environment - start script sources voicemode.env for port/model config
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/local/cuda/bin"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voicemode-whisper

[Install]
WantedBy=default.target