# packages/razer-control-center/razer_control_center-1.9.8.tar.gz/razer_control_center-1.9.8/packaging/systemd/razer-remap-daemon.service
[Unit]
Description=Razer Control Center Remap Daemon
Documentation=https://github.com/AreteDriver/Razer_Controls
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/razer-remap-daemon
Restart=on-failure
RestartSec=5

# Environment
Environment=PYTHONUNBUFFERED=1

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/razer-control-center

[Install]
WantedBy=default.target
